# raspi/systemd/ss-susi-voice@.service
[Unit]
Description=SUSI Linux
Wants=network-online.target
After=network-online.target ss-susi-server.service

[Service]
Type=simple
User=%i
ExecStart=susi-voice -v --short-log
Restart=on-failure

[Install]
WantedBy=multi-user.target

